# deploy/qiaichemy-backend.service
[Unit]
Description=QiAIchemy Backend API
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/srv/qiaichemy-backend
Environment=HOME=/home/ubuntu
EnvironmentFile=/srv/qiaichemy-backend/.env
Environment=NODE_ENV=production
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
ExecStart=/bin/bash -lc '. /home/ubuntu/.nvm/nvm.sh && nvm use 22.13.0 >/dev/null && exec node dist/server.js'
Restart=always
RestartSec=5
KillSignal=SIGTERM
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
